# conf/etc/s6/services/02-nginx/run
#!/bin/sh
service-handler "PHP-FPM" "nginx -g 'daemon off;'"