# conf/etc/s6/services/01-phpfpm/run
#!/bin/sh
NAME="PHP-FPM"
COMMAND="php-fpm7 -F"
WAIT="/tmp/01-nextcloud"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
