# conf/etc/periodic/15min/01-nextcloud
#!/bin/sh
su - nginx -s /bin/ash -c "php -f /var/www/nextcloud/cron.php"